# DEBIAN/postinst
#!/bin/bash

#Configuring the chonosDDNS to local.bot.chon.group
/usr/bin/chonosDDNSManager --conf -d local -u none -t none

systemctl daemon-reload
systemctl enable chonos-log.service
systemctl start  chonos-log.service
systemctl enable chonos-sysconfig.service
systemctl start  chonos-sysconfig.service
systemctl enable chonos-task.service
systemctl start  chonos-task.service
systemctl enable chonos-neighbors.service
systemctl start chonos-neighbors.service